# services/crosscompute-pages.service
[Service]
Type=oneshot
ExecStart=/root/Projects/landing-pages/scripts/build.sh
